# crontab
0 12 * * * python /Projects/weather/manage.py send_email
